# claude/templates/cron/com.user.claude.learn-eval.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude.learn-eval</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <!-- Replace USERNAME with your macOS short user name when installing, or use the installer's --with-learn-eval-cron flag which templates this automatically. -->
        <string>/Users/USERNAME/.claude/cron/learn-eval-monthly.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/.claude/learnings/review-evals/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/.claude/learnings/review-evals/logs/launchd.err</string>
</dict>
</plist>
